# etc/systemd/system/mirror-sync.service
[Unit]
Description=Sync Linux distro mirrors via rsync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mirror-sync.py

