# etc/systemd/system/mirror-index.service
[Unit]
Description=Generate UCSB repo mirror web index

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-mirror-index.py

