# etc/systemd/system/mirror-sync.timer
[Unit]
Description=Run mirror sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target


